# etc/arimport.cron
#!/bin/sh
# Run import for all configured apparelrow providers

su -s /bin/sh -c "/home/hanssonlarsson/apparelrow/bin/import" www-data
