# etc/dashboard_payment.cron
#!/bin/sh

su -s /bin/sh -c ". %(path)s/bin/activate && cd %(path)s/releases/current && ./manage.py dashboard_payment" %(run_user)s
