# etc/dashboard_import.cron
#!/bin/sh

su -s /bin/sh -c ". %(path)s/bin/activate && cd %(path)s/releases/current/apparelrow && ./manage.py dashboard_import --days 365" %(run_user)s
