# etc/import.cron
@daily cd /opt/hanssonlarsson/apparelrow; bin/import

