# etc/vendor_check.cron
#!/bin/sh

su -s /bin/sh -c ". %(path)s/bin/activate && cd %(path)s/releases/current && ./manage.py vendor_check" %(run_user)s
